# server/votes/cbw-votes-backup.timer
[Unit]
Description=Daily backup of CryptoBonusWorld voting data

[Timer]
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
